# kiro-worker.service
[Unit]
Description=Semantic API Kiro Worker
After=network.target semantic-api.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/semantic-api
Environment="AWS_REGION=us-east-1"
Environment="PATH=/home/ec2-user/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/bin/bash scripts/kiro-worker.sh
Restart=always
RestartSec=10
StandardOutput=append:/home/ec2-user/semantic-api/kiro-worker.log
StandardError=append:/home/ec2-user/semantic-api/kiro-worker.log

[Install]
WantedBy=multi-user.target
