# systemd/kiro-api.service
[Unit]
Description=Kiro API Server
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/kiro-api
Environment="NODE_ENV=production"
Environment="KIRO_API_PORT=8081"
Environment="AWS_REGION=us-east-1"
Environment="KIRO_QUEUE_TABLE=kiro-task-queue"
Environment="KIRO_SESSIONS_TABLE=kiro-sessions"
Environment="KIRO_MAX_SESSIONS=5"
ExecStart=/usr/bin/node src/kiro-api-server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kiro-api

[Install]
WantedBy=multi-user.target
